# supervisor.conf
[program:tv3]
command = python3 tv3.py
directory = /home/pi
user = pi
stdout_logfile = /var/log/tv.log
redirect_stderr = true
autostart = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8
